# examples/jsonmon.service
# jsonmon service for systemd

[Unit]
After      = network.target

[Service]
User             = jsonmon
WorkingDirectory = /etc/jsonmon
ExecStart        = /usr/local/bin/jsonmon config.yml
Restart          = always
LimitNOFILE      = 10240

# The JSON API host & port:
#Environment = HOST=localhost
#Environment = PORT=3000

[Install]
WantedBy   = multi-user.target
