# examples/jsonmon.service
[Unit]

[Service]
User             = jsonmon
ExecStart        = /usr/local/bin/jsonmon /etc/jsonmon/config.yml
Restart          = always

# The JSON API host & port:
#Environment=HOST=localhost
#Environment=PORT=3000

[Install]
WantedBy  = default.target
